# config/upstart_unicorn.conf.erb
description "Starts unicorn <%= @project_name %>"

start on stopped rc RUNLEVEL=[2345]
stop on starting runlevel [016]

chdir /u/apps/<%= @project_name %>/current
env RBENV_ROOT="/usr/local/rbenv"

console owner

script

PATH=/usr/local/rbenv/bin:/usr/local/rbenv/shims:$PATH

rbenv exec bundle exec unicorn_rails -c /u/apps/<%= @project_name %>/current/config/unicorn.rb -E <%= @stage %>

end script

respawn
